# sdk-setup/systemd/sdk-webappstub.socket
[Unit]
Description=Webapp Socket

[Socket]
ListenStream=9292
Accept=yes
BindIPv6Only=both

[Install]
WantedBy=sockets.target
