# sdk-setup/systemd/sfdk-bus-expose-nonce.service
[Unit]
Description=Expose D-Bus nonce on a location known to sfdk.

[Service]
ExecStart=/usr/libexec/sdk-setup/sfdk-bus-expose-nonce
